#  program_vm1.conf
[program:vm1]
command=/usr/sbin/firecracker --id vm1 --config-file /root/vm1.conf.json --api-sock /tmp/vm1.sock
stdout_logfile=/var/log/vm1.log
stderr_logfile=/var/log/vm1.error.log
autostart=true
